# packaging/debian/DEBIAN/postinst
#!/bin/bash
# postinst — Runs after dpkg -i install
set -e

APP_DIR="/usr/lib/sonic-openconfig-diagnostic-gateway"
VENV_DIR="${APP_DIR}/venv"
AUDIT_DIR="/var/log/sonic-openconfig-diagnostic-gateway"
ENV_FILE="/etc/sonic-openconfig-diagnostic-gateway/gateway.env"
SERVICE="sonic-openconfig-diagnostic-gateway"

echo "=== SONiC OpenConfig Diagnostic Gateway — Post-install ==="
echo ""

# ── 1. Create audit log directory ─────────────────────────────────
echo -n "Creating audit log directory... "
mkdir -p "${AUDIT_DIR}"
chmod 755 "${AUDIT_DIR}"
echo "${AUDIT_DIR}"

# ── 2. Secure the environment file ────────────────────────────────
echo -n "Securing environment file... "
chmod 600 "${ENV_FILE}" 2>/dev/null || true
chown root:root "${ENV_FILE}" 2>/dev/null || true
echo "done"

# ── 3. Create Python virtual environment ──────────────────────────
echo "Creating Python virtual environment (this may take ~60s)..."
if [ -f "${VENV_DIR}/bin/activate" ]; then
    echo "  Virtual environment already exists — skipping creation."
else
    # Remove stale broken venv from previous failed install
    rm -rf "${VENV_DIR}"
    # SONiC may not have python3-venv or ensurepip available.
    # Create venv without pip, then bootstrap pip manually.
    if python3 -m venv --without-pip "${VENV_DIR}" 2>/dev/null; then
        echo "  venv created (--without-pip)"
        # Bootstrap pip using ensurepip or get-pip.py
        "${VENV_DIR}/bin/python3" -m ensurepip --upgrade 2>/dev/null || {
            echo "  ensurepip not available — downloading get-pip.py..."
            curl -sS https://bootstrap.pypa.io/get-pip.py | "${VENV_DIR}/bin/python3"
        }
    else
        echo "  ERROR: Could not create virtual environment."
        echo "  Try manually: python3 -m venv --without-pip ${VENV_DIR}"
        exit 1
    fi
    echo "  venv ready at ${VENV_DIR}"
fi

# ── 4. Install Python dependencies ────────────────────────────────
echo "Installing Python dependencies..."
source "${VENV_DIR}/bin/activate"
pip install --upgrade pip -q 2>/dev/null || true

# Install from pyproject.toml (excludes gnmi/gnoi optional deps since
# we use gnmic/grpcurl CLI tools instead)
if [ -f "${APP_DIR}/pyproject.toml" ]; then
    pip install "${APP_DIR}/" 2>&1 | tail -5
else
    echo "  WARNING: pyproject.toml not found at ${APP_DIR}"
fi
deactivate

# ── 5. Register with systemd ──────────────────────────────────────
echo "Registering systemd service..."
systemctl daemon-reload
systemctl enable "${SERVICE}.service" 2>/dev/null || true

# ── 6. Print configuration instructions ───────────────────────────
echo ""
echo "============================================================"
echo " Installation complete!"
echo "============================================================"
echo ""
echo "NEXT STEPS:"
echo ""
echo "  1. Edit the environment file with your sandbox credentials:"
echo "     sudo vi ${ENV_FILE}"
echo ""
echo "  2. Start the service:"
echo "     sudo systemctl start ${SERVICE}"
echo ""
echo "  3. Check status:"
echo "     sudo systemctl status ${SERVICE}"
echo ""
echo "  4. Verify:"
echo "     curl http://localhost:8080/health"
echo ""
echo "  View logs:"
echo "     sudo journalctl -u ${SERVICE} -f"
echo ""
